# .config/systemd/user/set-random-background.timer
[Unit]
Description=Set random background timer

[Timer]
OnActiveSec=0
OnUnitActiveSec=10m

[Install]
WantedBy=timers.target
