# .config/systemd/user/set-random-background.service
[Unit]
Description=Set random background

[Service]
Type=oneshot
ExecStart=/home/trygve/bin/set_random_background.pl
