# etc/systemd/system/dhcpcd.service.d/execstart.conf
[Service]
ExecStart=
ExecStart=/usr/bin/dhcpcd -q
