# etc/systemd/user/update-net-widget.service
[Unit]
Description=Update net widget

[Service]
ExecStart=/home/trygve/dev/update-net-widget/update-net-widget

[Install]
WantedBy=default.target
