# etc/udev/rules.d/90-keyboard-connect.rules
SUBSYSTEM=="input", ACTION=="add", RUN+="/bin/su trygve -c 'export DISPLAY=:0 && /usr/bin/xmodmap /home/trygve/.Xmodmap'"
