# etc/udev/rules.d/90-monitor-connect.rules
SUBSYSTEM=="drm", ACTION=="change", RUN+="/bin/su trygve -c '/home/trygve/bin/update-monitor'"
